# debian/ws-camerad.service
[Unit]
Description=High-Performance Camera Daemon
Documentation=man:ws-camerad(8)
After=network.target ws-audiod.service

[Service]
Type=simple
ExecStart=/usr/bin/ws-camerad -c /etc/ws/camerad/ws-camerad.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
